# pivccu/rtc/rpi-rf-mod_systohc.timer
[Unit]
 
[Timer]
OnUnitActiveSec=1h
OnBootSec=5m
 
[Install]
WantedBy=timers.target
