# package/pivccu/prerm
#!/bin/sh
systemctl stop pivccu.service
systemctl disable pivccu.service

rm -f /usr/sbin/pivccu-attach
rm -f /usr/sbin/pivccu-info
rm -f /usr/sbin/pivccu-device
