# package/pivccu-modules-raspberrypi/postinst
#!/bin/bash
. /usr/share/debconf/confmodule

set -e

case "$1" in
  configure)
    if [ -f /usr/share/rpikernelhack/overlays/pivccu-raspberrypi.dtbo ]; then
      rm -f /var/lib/piVCCU/dtb/overlays/pivccu-raspberrypi.dtbo
      dpkg-divert --package rpikernelhack --remove --rename /var/lib/piVCCU/dtb/overlays/pivccu-raspberrypi.dtbo
      sync
    fi

    sed -i /var/lib/piVCCU/dtb/config.txt -e '/dtoverlay=pivccu-raspberrypi/d'
    [ -n "$(tail -1 /var/lib/piVCCU/dtb/config.txt)" ] && echo >> /var/lib/piVCCU/dtb/config.txt
    echo "dtoverlay=pivccu-raspberrypi" >> /var/lib/piVCCU/dtb/config.txt

    umount /var/lib/piVCCU/dtb

    systemctl enable pivccu-rpi-modules.service

    if [ ! -e /pivccu.createimage ]; then
      db_fset pivccu/reboot_required seen false
      db_input medium pivccu/reboot_required || true
      db_go
    fi
    ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    ;;
esac

